# db/crontab
0 * * * * neo4j /scripts/backup.sh >> /var/log/cron/backup.log 2>&1
